# service/kotlin-app-template.service
[Unit]
Description=kotlin-app-template

[Service]
User=deployer
WorkingDirectory=/services/kotlin-app-template
ExecStart=/services/kotlin-app-template/run
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target